# debian/postinst
#!/bin/sh
# postinst script for debianissimo-utils.
#
# See: dh_installdeb(1).

set -e

# Summary of how this script can be called:
#        * <postinst> 'configure' <most-recently-configured-version>
#        * <old-postinst> 'abort-upgrade' <new version>
#        * <conflictor's-postinst> 'abort-remove' 'in-favour' <package>
#          <new-version>
#        * <postinst> 'abort-remove'
#        * <deconfigured's-postinst> 'abort-deconfigure' 'in-favour'
#          <failed-install-package> <version> 'removing'
#          <conflicting-package> <version>
# for details, see https://www.debian.org/doc/debian-policy/ or
# the debian-policy package.

FSTAB="/etc/fstab"
FSTAB_BACKUP="/etc/fstab.debianissimo-utils-backup"

LABEL="Memorissimo"
DISK="LABEL=$LABEL"
MOUNTPOINT="/media/$LABEL"
FSTYPE="auto"
OPTIONS="defaults,rw,nofail,x-systemd.automount"
OPTIONS="$OPTIONS,x-systemd.requires=fsck@dev-disk-by\\x2dlabel-$LABEL.service"

case "$1" in
    configure)
        if [ -f "$FSTAB" ] && ! grep -q "^$DISK\s" "$FSTAB"; then
            cp -n "$FSTAB" "$FSTAB_BACKUP"
            {
                echo "# Automount entry for $LABEL"
                echo "# Based on an original Ordissimo configuration, modified for Debianissimo."
                echo "# Debianissimo is not affiliated with Ordissimo in any way."

                echo "$DISK $MOUNTPOINT $FSTYPE $OPTIONS 0 0"
            } >> "$FSTAB" || cp "$FSTAB_BACKUP" "$FSTAB"

            if [ ! -d "$MOUNTPOINT" ]; then
                mkdir -p "$MOUNTPOINT"
            fi

            if [ -f "$FSTAB_BACKUP" ]; then
                rm "$FSTAB_BACKUP"
            fi
        fi
        if [ -x /usr/bin/systemctl ]; then
            systemctl enable update-sysinfo.service || true
        fi
        

    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument '$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
